# packaging/services/launchd/Library/LaunchDaemons/sensu-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Disabled</key><false/>
    <key>Label</key><string>sensu-agent</string>
    <key>Program</key><string>/usr/bin/sensu-agent</string>
    <key>ProgramArguments</key><array>
      <string>/usr/bin/sensu-agent</string>
      <string>start</string>
    </array>
    <key>KeepAlive</key> <true/>
    <key>StandardOutPath</key> <string>/var/log/sensu-agent-stdout.log</string>
    <key>StandardErrorPath</key> <string>/var/log/sensu-agent-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
    </dict>
    <!-- Things to implement eventually.
    <key>UserName</key><string>sensu</string>
    <key>GroupName</key><string>sensu</string>
    <key>Umask</key><integer>022</integer>
    <key>RootDirectory</key><string>/</string>
    <key>WorkingDirectory</key><string>/</string>
    <key>Nice</key><integer>0</integer>
    <key>SoftResourceLimits</key><dict>
      <key>Core</key><integer></integer>
      <key>CPU</key><integer></integer>
      <key>Data</key><integer></integer>
      <key>FileSize</key><integer></integer>
      <key>MemoryLock</key><integer></integer>
      <key>NumberOfFiles</key><integer></integer>
      <key>NumberOfProcesses</key><integer></integer>
      <key>ResidentSetSize</key><integer></integer>
      <key>Stack</key><integer></integer>
    </dict>
    <key>HardResourceLimits</key><dict>
      <key>Core</key><integer></integer>
      <key>CPU</key><integer></integer>
      <key>Data</key><integer></integer>
      <key>FileSize</key><integer></integer>
      <key>MemoryLock</key><integer></integer>
      <key>NumberOfFiles</key><integer></integer>
      <key>NumberOfProcesses</key><integer></integer>
      <key>ResidentSetSize</key><integer></integer>
      <key>Stack</key><integer></integer>
    </dict>
      -->
  </dict>
</plist>

